# rsyslog/supervisor/rsyslog.ini
[program:rsyslog]
command=/usr/local/sbin/rsyslogd -n
priority=1

